# tools/exporters/com.agentic.mac-gpu-exporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic.mac-gpu-exporter</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jiwoolee/Project/agentic-ai/tools/exporters/mac_gpu_exporter.py</string>
        <string>--port</string>
        <string>9101</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/mac-gpu-exporter.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mac-gpu-exporter.err</string>
</dict>
</plist>
